# etc/root/systemd/system/getty@tty1.service.d/override.conf
[Service]  
ExecStart=  
ExecStart=-/sbin/agetty --autologin vn36 --noclear %I 38400 linux