# lib/dotfiles/crontab
*/5 * * * * /home/vn36/.local/bin/scripts/bat_alert.sh
# */10 * * * * /home/vn36/.local/bin/scripts/bat_log.sh
* */1 * * * anacron -t "$HOME"/etc/anacrontab -S "$HOME"/var/spool/anacron

