# .local/etc/crontab
*/10 * * * * /home/vn36/.local/bin/custom/bat_alert.sh
*/10 * * * * /home/vn36/.local/bin/custom/bat_log.sh

